# systemd/care-cat-status.service
[Unit]
Description=CARECatStatus - Real-time cat status board
After=network.target

[Service]
Type=simple
User=USER_PLACEHOLDER
WorkingDirectory=REPO_DIR_PLACEHOLDER
Environment="PORT=3007"
Environment="DATABASE_URL=REPO_DIR_PLACEHOLDER/cats.db"
Environment="RUST_LOG=info"
ExecStart=REPO_DIR_PLACEHOLDER/target/release/care-cat-status
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
